# package/network/utils/bwm/files/bwm.init
#!/bin/sh /etc/rc.common

START=99
STOP=99

start() {
	wan_ifname=$(uci get network.wan.ifname)
	bwm -i $wan_ifname -t 2 &
}

stop() {
	killall -9 bwm
}
