# package/base-files-qca/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

/etc/init.d/zyreload.sh &
exit 0
