# target/linux/brcm963xx/base-files-416L05/etc/inittab
# Default entries for busybox init.
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::restart:/sbin/init
#::askfirst:-/bin/sh
::sysinit:/etc/init.d/rcS
#::respawn:-/bin/sh
::askfirst:/sbin/getty -L ttyS0 115200 vt100
